# deploy/systemd/tapdash-pipeline-daemon.service
[Unit]
Description=TapDash Pipeline Daemon
After=network.target

[Service]
Type=simple
User=tapdash
Group=tapdash
WorkingDirectory=/opt/tapdash-swarm
EnvironmentFile=/etc/tapdash/swarm.env
ExecStart=/usr/bin/python3 /opt/tapdash-swarm/pipeline_daemon.py --interval-seconds 10
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
